# deploy/systemd/api-quiz-bank-protected-beta-telegram.service.d/10-runtime-user-and-assets.conf
[Service]
ExecStartPre=/usr/bin/install -d -o 65532 -g 65532 -m 0750 /var/lib/api-quiz-bank/app-var
ExecStartPre=/usr/bin/install -d -o 65532 -g 65532 -m 0750 /var/lib/api-quiz-bank/app-var/visual-assets
ExecStart=
ExecStart=/bin/sh -c '/usr/bin/docker exec --user 65532:65532 api-quiz-bank-pilot python /app/tools/run_protected_beta_schedule.py --mode real --approve-real-send --token-file /run/secrets/api_quiz_bank_telegram_bot_token --slot-time "$(TZ=Europe/Berlin date +%%H:%%M)"'
